# systemd/xpl-domogeek.service
[Unit]
Description=Wiseflat xPL domogeek daemon
After=networking.service
After=xpl-hub.service

[Service]
Type=simple
ExecStart=node /opt/xpl-domogeek/index.js
Restart=always

[Install]
WantedBy=multi-user.target
